# tools/build-system/launchpad/sflphone-daemon-video/debian/cron.d
#
# Regular cron jobs for the sflphone package
#
0 4	* * *	root	sflphone_maintenance
